# squashfs-root/etc/uci-defaults/luci-dev-feature
#!/bin/sh

cp /rom/etc/config/dev_feature /etc/config/dev_feature
