# squashfs-root/etc/init.d/kms
#!/bin/sh /etc/rc.common

START=50

start() {
/etc/vlmcsd/vlmcsd -i /etc/vlmcsd/vlmcsd.ini -p /tmp/vlmcsd.pid -l /tmp/vlmcsd.log
}
stop() {
 killall -9 vlmcsd
}
