# squashfs-root/etc/hotplug.d/iface/75-remote
#!/bin/sh

if [ "$ACTION" = "ifup" ]; then
	if [ "$INTERFACE" = "wan" -o "$INTERFACE" = "vpn" ]; then
		/etc/init.d/remote restart
	fi
fi
