# squashfs-root/usr/lib/opkg/info/dnscrypt-proxy.list
/etc/config/dnscrypt-proxy
/usr/share/dnscrypt-proxy/dnscrypt-resolvers.csv
/usr/sbin/dnscrypt-proxy
/etc/init.d/dnscrypt-proxy
