# squashfs-root/usr/lib/opkg/info/http.list
/usr/sbin/http
/etc/init.d/http
/etc/config/http_report
