# squashfs-root/usr/lib/opkg/info/luci-mod-base.list
/usr/lib/lua/luci/view/h5/wifiConfig.htm
/usr/lib/lua/luci/controller/admin/page.lua
/usr/lib/lua/luci/view/userAgreementContent-ru-ru.htm
/www/luci-static/images/topology.png
/etc/config/luci
/usr/lib/lua/luci/phicomm/validator.lua
/usr/lib/lua/luci/data/wireless_plt.lua
/usr/lib/lua/luci/controller/admin/scheme_upgrade_plt.lua
/www/luci-static/images/button_disable.png
/www/luci-static/images/signal_4.png
/usr/lib/lua/luci/controller/admin/wireless_plt.lua
/www/luci-static/images/check_upgrade_success.png
/www/luci-static/images/wisp_h5.png
/etc/init.d/luci_fixtime
/usr/lib/lua/luci/view/pc/netState.htm
/usr/lib/lua/luci/controller/admin/protest.lua
/www/luci-static/images/router.svg
/www/luci-static/images/phicomm.png
/www/luci-static/images/signal.gif
/www/luci-static/dynaform/up.css
/www/luci-static/images/pic_code.png
/www/luci-static/images/signal_lock_1.png
/usr/lib/lua/luci/view/h5/setLgPwd.htm
/usr/lib/lua/luci/view/language.htm
/usr/lib/lua/luci/phicomm/setting.lua
/www/luci-static/images/wifi.png
/www/luci-static/dynaform/icon_fonts_style.css
/www/luci-static/images/wisp_pc.png
/www/luci-static/dynaform/frame.css
/www/luci-static/dynaform/h5-en-us.css
/www/luci-static/images/color_line.png
/www/luci-static/images/upgrade_rocket.png
/usr/lib/lua/luci/view/footer.htm
/usr/lib/lua/luci/view/h5/guide.htm
/usr/lib/lua/luci/data/time_reboot_plt.lua
/usr/lib/lua/luci/view/error404.htm
/usr/lib/lua/luci/controller/admin/dev_feature.lua
/www/luci-static/dynaform/feature.js
/www/luci-static/dynaform/pc-th-th.css
/www/luci-static/images/signal_lock_2.png
/usr/lib/lua/luci/view/h5/wifiList.htm
/www/luci-static/images/check_button.png
/usr/lib/lua/luci/controller/ds.lua
/www/luci-static/fonts/iconfont.woff
/usr/lib/lua/luci/view/h5/menu.htm
/www/luci-static/images/game_bg.png
/usr/lib/lua/luci/data/guide.lua
/www/luci-static/dynaform/upgrade.js
/www/luci-static/images/dog.png
/www/luci-static/images/close_time.png
/www/luci-static/lib/json.min.js
/www/luci-static/images/iptv_port.png
/usr/lib/lua/luci/view/pc/guideWifiSet.htm
/usr/lib/lua/luci/phicomm/debug.lua
/usr/lib/lua/luci/view/h5/guideWifiSet.htm
/www/luci-static/images/signal_ok.png
/usr/lib/lua/luci/view/h5/timeReboot.htm
/www/luci-static/dynaform/it-it.js
/www/luci-static/dynaform/routing.js
/www/luci-static/images/network.png
/www/luci-static/images/pic_light.png
/www/luci-static/dynaform/pc-en-us.css
/www/luci-static/images/pic_code_de.png
/www/luci-static/dynaform/de-de.js
/www/luci-static/lib/libs.js
/www/luci-static/images/close_wifi.png
/www/luci-static/images/router_blue.png
/usr/lib/lua/luci/view/pc/setLgPwd.htm
/www/luci-static/images/close_list.png
/www/luci-static/dynaform/pc-ko-kr.css
/usr/lib/lua/luci/controller/admin/time_reboot.lua
/www/luci-static/dynaform/pc-id-id.css
/www/luci-static/lib/flexible.min.js
/www/luci-static/images/signal_lock_3.png
/usr/lib/lua/luci/view/pc/userAgreement.htm
/etc/cfgAdaptor/system.sh
/usr/lib/lua/luci/controller/admin/wireless.lua
/www/luci-static/dynaform/vendor.js
/usr/lib/lua/luci/view/pc/guide.htm
/usr/lib/lua/luci/data/time_reboot.lua
/www/luci-static/dynaform/pc-vi-vi.css
/www/luci-static/images/pic_router.png
/www/luci-static/images/pic_left.png
/usr/lib/lua/luci/view/footer-welcome.htm
/www/luci-static/dynaform/pc-uk-ua.css
/www/luci-static/images/icon_arrow.png
/usr/lib/lua/luci/controller/admin/device_plt.lua
/usr/lib/lua/luci/controller/admin/time_reboot_plt.lua
/usr/lib/lua/luci/view/userAgreementContent-en-us.htm
/www/luci-static/dynaform/verify.js
/usr/lib/lua/luci/data/device_mng_plt.lua
/www/luci-static/dynaform/zh-tw.js
/www/luci-static/dynaform/h5.js
/www/luci-static/dynaform/zh-cn.js
/usr/lib/lua/luci/phicomm/util.lua
/usr/lib/lua/luci/controller/admin/device.lua
/www/luci-static/dynaform/error.js
/usr/lib/lua/luci/view/wifi-config-script.htm
/www/luci-static/dynaform/h5-de-de.css
/www/luci-static/images/signal_lock_4.png
/www/luci-static/images/uncheck_button.png
/www/luci-static/images/button.png
/usr/lib/lua/luci/phicomm/statistic.lua
/www/luci-static/lib/stateman.min.js
/www/luci-static/dynaform/default_feature.js
/www/luci-static/images/circle.png
/etc/config/diagnosis
/etc/uci-defaults/luci-dev-feature
/www/luci-static/dynaform/cityList.js
/usr/lib/lua/luci/view/userAgreementContent-uk-ua.htm
/www/luci-static/lib/respond.min.js
/www/luci-static/dynaform/help.js
/usr/lib/lua/luci/controller/admin/device_mng.lua
/www/luci-static/dynaform/frame.js
/www/luci-static/images/app-icon/reboot.png
/www/luci-static/images/speed_bg.png
/usr/lib/lua/luci/phicomm/fs.lua
/etc/config/timereboot
/etc/init.d/luci_generatecache
/etc/cfgAdaptor/luci.sh
/www/luci-static/images/loading.gif
/usr/lib/lua/luci/phicomm/statistic_plt.lua
/usr/lib/lua/luci/data/wireless.lua
/usr/lib/lua/luci/view/userAgreementContent.htm
/www/luci-static/images/password.png
/usr/lib/lua/luci/controller/admin/scheme_upgrade.lua
/usr/lib/lua/luci/view/index.htm
/usr/lib/lua/luci/controller/admin/device_mng_plt.lua
/www/luci-static/dynaform/pc-ja-jp.css
/usr/lib/lua/luci/view/pc/wifiConfig.htm
/www/luci-static/dynaform/pc-zh-tw.css
/usr/lib/lua/luci/data/device_mng.lua
/www/luci-static/dynaform/h5.css
/www/luci-static/images/router-white.png
/usr/lib/lua/luci/dispatcher.lua
/www/luci-static/fonts/iconfont.eot
/www/luci-static/dynaform/es-es.js
/www/luci-static/images/router.png
/www/luci-static/dynaform/pc-de-de.css
/usr/lib/lua/luci/view/h5/editClient.htm
/www/luci-static/images/net.svg
/www/luci-static/dynaform/pc.css
/usr/lib/lua/luci/view/h5/userAgreement.htm
/usr/lib/lua/luci/view/h5/login.htm
/www/luci-static/dynaform/up.js
/www/luci-static/lib/jquery.min.js
/www/luci-static/dynaform/vi-vi.js
/www/luci-static/images/signal_1.png
/www/luci-static/dynaform/ja-jp.js
/usr/lib/lua/luci/data/guide_plt.lua
/www/luci-static/dynaform/ru-ru.js
/www/luci-static/images/device_empty.png
/www/luci-static/dynaform/fr-fr.js
/www/luci-static/dynaform/th-th.js
/www/luci-static/dynaform/en-us.js
/www/luci-static/images/list_empty.png
/www/luci-static/images/pic_computer.png
/www/luci-static/images/game_success.png
/www/luci-static/dynaform/pt-pt.js
/usr/lib/lua/luci/view/pc/menu.htm
/www/luci-static/images/computer.svg
/usr/lib/lua/luci/view/pc/deviceManage.htm
/www/luci-static/fonts/iconfont.ttf
/www/luci-static/fonts/iconfont.svg
/www/luci-static/dynaform/pc-pt-pt.css
/usr/lib/lua/luci/view/h5/deviceManage.htm
/www/luci-static/images/signal_bg.png
/www/luci-static/images/phiwifi.png
/www/luci-static/images/pic_code_us.png
/www/luci-static/lib/ajax.js
/usr/lib/lua/luci/view/userAgreementContent-de-de.htm
/www/luci-static/images/bg_map.jpg
/www/luci-static/dynaform/id-id.js
/www/luci-static/images/check_upgrade_failed.png
/usr/lib/lua/luci/phicomm/lib/shell.lua
/www/luci-static/images/game_ok.png
/etc/config/dev_feature
/usr/lib/lua/luci/view/pc/App.htm
/www/luci-static/images/map_bg.png
/www/luci-static/images/signal_2.png
/www/luci-static/dynaform/pc.js
/www/luci-static/images/pic_light_blue.png
/www/luci-static/dynaform/ie.css
/usr/lib/lua/luci/view/pc/upgradeWgt.htm
/www/luci-static/lib/echarts.simple.min.js
/usr/lib/lua/luci/view/h5/netState.htm
/usr/lib/lua/luci/phicomm/error.lua
/usr/lib/lua/luci/view/pc/login.htm
/www/luci-static/dynaform/ko-kr.js
/www/luci-static/images/diagnose_h5.png
/www/luci-static/images/menu_code.png
/www/luci-static/dynaform/pc-it-it.css
/usr/lib/lua/luci/view/pc/modifyPwd.htm
/usr/lib/lua/luci/view/error500.htm
/www/luci-static/dynaform/product.js
/www/luci-static/images/game_empty.png
/www/luci-static/dynaform/pc-ru-ru.css
/www/luci-static/dynaform/pc-es-es.css
/www/luci-static/images/router_wisp.svg
/usr/lib/lua/luci/phicomm/lib/bit.lua
/usr/lib/lua/luci/view/pc/timeReboot.htm
/www/luci-static/images/bg_botton.png
/www/luci-static/images/app-icon/lock.png
/www/luci-static/images/signal_3.png
/www/luci-static/images/favicon.ico
/usr/lib/lua/luci/view/h5/modifyPwd.htm
/www/luci-static/dynaform/uk-ua.js
/www/luci-static/dynaform/pc-fr-fr.css
/www/luci-static/images/pic_bottombar.png
/usr/lib/lua/luci/view/h5/App.htm
