# squashfs-root/usr/lib/opkg/info/sleep_timer.list
/etc/init.d/sleep_timer
/etc/config/sleep_timer
/usr/bin/sleep_timer
